# recipe_book/src/chef_recipes/recipes/scheduled-task-persistence/recipe.yml
name: scheduled-task-persistence
version: "1.0"
description: >
  Emulate persistence via scheduled task creation (T1053.005). Validates
  whether endpoint detection fires when schtasks.exe is used to create
  tasks pointing to suspicious file paths, a common persistence mechanism.

metadata:
  author: scott-thornton
  mitre_techniques:
    - T1053.005
  mitre_tactics:
    - persistence
    - execution
  difficulty: beginner
  estimated_time: 10m
  tags:
    - windows
    - persistence
    - execution
    - scheduled-task
    - purple-team
    - purplechef

mise_en_place:
  terraform_module: windows-target
  ansible_roles:
    - limacharlie-sensor
    - atomic-install
  prerequisites:
    caldera_server_running: true
    limacharlie_org_configured: true

attack:
  method: caldera
  caldera:
    adversary_name: chef-scheduled-task-persistence
    abilities:
      - technique_id: T1053.005
        ability_id: auto
        description: Create a scheduled task pointing to a suspicious path for persistence
    group: chef-targets
    timeout: 300

validate:
  detection_source: limacharlie
  wait_seconds: 90
  expected_rules:
    - name: Scheduled Task Creation
      tags:
        - t1053
  sigma_rules:
    - path: sigma-rules/schtasks-creation.yml

report:
  format:
    - json
    - html
  evidence_required: true

advise:
  generate_sigma: true
  generate_kql: false
